# root/etc/s6-overlay/s6-rc.d/init-broker-launcher/up
#!/usr/bin/with-contenv bash

echo "[broker-mod] Installing pcsx2-launcher.sh..."
cp /root/pcsx2-launcher.sh /config/pcsx2-launcher.sh
chmod +x /config/pcsx2-launcher.sh
chown abc:abc /config/pcsx2-launcher.sh

echo "[broker-mod] Patching labwc autostart to use launcher..."
AUTOSTART="/config/.config/labwc/autostart"

# Only patch if not already patched
if grep -q "pcsx2-launcher.sh" "$AUTOSTART" 2>/dev/null; then
    echo "[broker-mod] Autostart already patched, skipping"
else
    # Replace the pcsx2-qt line with our launcher
    sed -i 's|foot -e /usr/bin/pcsx2-qt|foot -e /config/pcsx2-launcher.sh|g' "$AUTOSTART"
    sed -i 's|xterm -e /usr/bin/pcsx2-qt|xterm -e /config/pcsx2-launcher.sh|g' "$AUTOSTART"
    echo "[broker-mod] Autostart patched successfully"
fi