# root/etc/s6-overlay/s6-rc.d/init-broker-launcher/up
/etc/s6-overlay/s6-rc.d/init-broker-launcher/init.sh