# root/etc/s6-overlay/s6-rc.d/init-pcsx2-config/up
/etc/s6-overlay/s6-rc.d/init-pcsx2-config/init.sh