# root/etc/s6-overlay/s6-rc.d/svc-broker/finish
#!/usr/bin/with-contenv bash
echo "[broker-mod] broker.py exited with code $1 — s6 will restart it"