# root/etc/s6-overlay/s6-rc.d/svc-broker/run
#!/usr/bin/with-contenv bash

echo "[broker-mod] Checking for Python3..."
if ! command -v python3 &> /dev/null; then
    echo "[broker-mod] Python3 not found. Installing..."
    apt-get update -qq && apt-get install -y -qq python3
else
    echo "[broker-mod] Python3 already installed: $(python3 --version)"
fi

echo "[broker-mod] Waiting for PCSX2.ini..."
INI="/config/.config/PCSX2/inis/PCSX2.ini"
for i in $(seq 1 60); do
    if [ -f "$INI" ]; then
        echo "[broker-mod] Patching PCSX2.ini..."

        # Enable PINE IPC
        sed -i 's/EnablePINE = false/EnablePINE = true/' "$INI"

        # Fix BIOS path to absolute container path
        sed -i 's|Bios = ../../bios|Bios = /config/bios|' "$INI"

        if ! grep -q "SDL-0" "$INI"; then
            sed -i '/^\[Pad1\]/a Type = DualShock2' "$INI"
        fi

        echo "[broker-mod] PCSX2.ini patched:"
        grep -E "EnablePINE|Bios " "$INI"
        break
    fi
    sleep 2
done

echo "[broker-mod] Setting broker.py and launcher executable..."
chmod +x /root/broker.py
chmod 755 /usr/local/bin/pcsx2-launcher.sh

echo "[broker-mod] Starting broker.py..."
exec python3 -u /root/broker.py