# root/etc/s6-overlay/s6-rc.d/svc-broker/run
#!/usr/bin/with-contenv bash
# Runs as root intentionally — needs to chmod/chown config dirs and manage system processes.

echo "[broker-mod] Checking for Python3..."
if ! command -v python3 &> /dev/null; then
    echo "[broker-mod] Python3 not found. Installing..."
    apt-get update -qq && apt-get install -y -qq python3
else
    echo "[broker-mod] Python3 already installed: $(python3 --version)"
fi

INI="/config/.config/PCSX2/inis/PCSX2.ini"
if [ -f "$INI" ]; then
    echo "[broker-mod] Existing config found — applying patches..."
else
    echo "[broker-mod] First run detected — waiting for setup wizard to complete (up to 60s)..."
fi

INI_FOUND=0
for i in $(seq 1 120); do
    if [ -f "$INI" ]; then
        INI_FOUND=1
        echo "[broker-mod] PCSX2.ini found after $((i / 2))s — patching..."

        # Enable PINE IPC
        sed -i 's/EnablePINE = false/EnablePINE = true/' "$INI"

        # Fix BIOS path to absolute container path
        sed -i 's|^Bios = .*|Bios = /config/bios|' "$INI"

        # PCSX2 writes relative paths here on exit; pin them to absolute
        sed -i 's|^Savestates = .*|Savestates = /config/.config/PCSX2/sstates|' "$INI"
        sed -i 's|^MemoryCards = .*|MemoryCards = /config/.config/PCSX2/memcards|' "$INI"

        echo "[broker-mod] PCSX2.ini patched:"
        grep -E "EnablePINE|Bios " "$INI"
        break
    fi
    # Log progress every 5 seconds
    if (( i % 10 == 0 )); then
        echo "[broker-mod] Still waiting for PCSX2.ini... ($((i / 2))s elapsed)"
    fi
    sleep 0.5
done
if (( INI_FOUND == 0 )); then
    echo "[broker-mod] WARNING: PCSX2.ini never appeared — PINE IPC will not be enabled until next restart"
fi

echo "[broker-mod] Setting broker.py and launcher executables..."
chmod +x /root/broker.py
chmod 755 /usr/local/bin/pcsx2-launcher.sh
chmod 755 /usr/local/bin/pcsx2-autostart.sh

# cleanup orphaned signal files from previous runs to prevent loops
if [ -f "/tmp/pcsx2-rom" ]; then
    echo "[broker-mod] Cleaning up orphaned signal file /tmp/pcsx2-rom..."
    rm -f /tmp/pcsx2-rom
fi

echo "[broker-mod] Writing labwc autostart..."
AUTOSTART="/config/.config/labwc/autostart"
mkdir -p "$(dirname "$AUTOSTART")"
# labwc runs autostart line-by-line in subshells -- multi-line bash cannot be used here.
printf '/usr/local/bin/pcsx2-autostart.sh\n' > "$AUTOSTART"
chmod +x "$AUTOSTART"

echo "[broker-mod] Starting broker.py..."
exec python3 -u /root/broker.py